# required.service
[Unit]
Description=__exists__
After=syslog.target network.target

[Service]
Restart=always

[Install]
WantedBy=multi-user.target

